# data/morning_call_latest.json
{
  "generated_at": "2026-01-25T06:49:39.340329+00:00",
  "analysis_date": "2026-01-25",
  "model": "gpt-5.1",
  "window_hours": 24,
  "source_file": "data/news_recent.json",
  "source_generated_at": "2026-01-25T06:00:17.173914+00:00",
  "source_days_back": 30,
  "source_total_items": 7202,
  "total_items_in_window_all": 89,
  "total_items_in_window_curated": 25,
  "morning_call_markdown": "### Executive Summary\n- CISA confirms active exploitation of a critical VMware vCenter RCE (CVE‑2024‑37079) – high risk for any exposed or unpatched virtualization infrastructure.\n- New ransomware and wiper activity (Osiris BYOVD ransomware, DynoWiper linked to Sandworm, SSO‑related data theft by ShinyHunters) underlines the need for stronger EDR defenses, identity monitoring, and backup/segmentation checks.\n- Targeted phishing (Marriott lookalike domains, multi‑stage campaigns delivering RAT + ransomware) continues to be the primary entry vector; email/identity telemetry and user reporting remain critical.\n\n---\n\n### High-priority items (immediate attention)\n\n- **Actively exploited VMware vCenter RCE (CVE‑2024‑37079 – CISA KEV)**\n  - What happened: CISA added a Broadcom VMware vCenter Server RCE to KEV and reports active exploitation in the wild.\n  - Why it matters operationally: Compromise of vCenter = full control of virtual infrastructure, lateral movement, and potential ransomware at scale.\n  - Recommended immediate actions for the SOC today:  \n    - Confirm with infra/VMware owners: inventory vCenter instances, internet exposure, and patch status; flag any unpatched or internet‑facing as “at risk.”  \n    - Increase monitoring for anomalous logins/API calls to vCenter (odd admin accounts, abnormal source IPs, new plugins/extensions).  \n    - If any vCenter is both exposed and unpatched, recommend emergency network restriction (VPN‑only / ACLs) and initiate proactive compromise assessment.\n\n- **Osiris ransomware using BYOVD to kill security tools**\n  - What happened: New Osiris ransomware variant brings its own vulnerable driver (BYOVD) to disable AV/EDR before encrypting.\n  - Why it matters operationally: Directly targets our core detection stack; infections may appear “silent” just before encryption.\n  - Recommended immediate actions for the SOC today:  \n    - Hunt for driver installation events involving unsigned/old/vulnerable kernel drivers and sudden EDR/AV service stoppage across endpoints/servers.  \n    - Set high‑priority alerts for: mass file renames, rapid creation of ransom notes, abnormal spike in 3xx/4xx file access errors on shares, combined with security tool shutdown.\n\n- **ShinyHunters SSO‑account data theft & general SSO targeting**\n  - What happened: ShinyHunters claim responsibility for theft of SSO‑related account data in multiple organizations.\n  - Why it matters operationally: Any IdP/SSO compromise can provide direct access into SaaS and internal apps, bypassing many traditional controls.\n  - Recommended immediate actions for the SOC today:  \n    - Review last 7–14 days of IdP logs for: impossible travel, unusual OAuth consent grants, suspicious new SAML/OIDC apps, MFA fatigue patterns.  \n    - Confirm conditional access/MFA enforced for all privileged and remote access; raise any gaps immediately to IAM/security engineering.\n\n- **Phishing using typo‑squatted domains & multi‑stage RAT→ransomware chains**\n  - What happened: Campaigns abuse lookalike domains (e.g., “rn” vs “m”) to impersonate brands like Marriott, plus a separate multi‑stage phishing campaign delivering Amnesia RAT and ransomware.\n  - Why it matters operationally: Direct user targeting with brand impersonation; foothold for RATs, data theft, and eventual encryption.\n  - Recommended immediate actions for the SOC today:  \n    - Tune email/security tools to flag lookalike domains (rn vs m, l vs I, etc.) and brand‑impersonation patterns; block any newly observed lookalike domains.  \n    - Hunt for recent email clicks leading to new processes: Office/PDF → script (VBS/JS/PowerShell) → network beacons → new services/scheduled tasks on user machines.\n\n- **DynoWiper / Sandworm OT‑adjacent destructive activity**\n  - What happened: New DynoWiper malware used in an attempted Sandworm attack against Poland’s power sector.\n  - Why it matters operationally: Shows active destructive tradecraft against critical infrastructure; similar TTPs could be reshaped to target financial operations (destructive attacks on Windows servers, AD, backup).\n  - Recommended immediate actions for the SOC today:  \n    - Validate backup/restore capability for AD, core banking, and file servers (including offline/immutable backups).  \n    - Hunt on critical servers for: bulk file deletion, suspicious use of `vssadmin`/`wbadmin`, and tools/scripts that overwrite MBR/boot records or wipe specific directories.\n\n---\n\n### Monitoring & detection recommendations\n\n- **VMware vCenter exploitation (CVE‑2024‑37079)**\n  - Log sources: vCenter/ESXi logs, management network firewalls, VPN concentrators, vulnerability scanner results.  \n  - Focus: anomalous admin logons, new or modified plugins/extensions, unexpected connections to vCenter from non‑admin networks. (MITRE: T1190, T1068, T1078)\n\n- **Ransomware/BYOVD & destructive behaviors (Osiris, DynoWiper)**\n  - Log sources: EDR, Windows Event Logs, Sysmon, file server logs.  \n  - Focus: new/unknown drivers installed (`Service Control Manager` events), termination of security processes, `vssadmin delete shadows`, `wmic shadowcopy delete`, spikes in file renames/deletes. (T1068, T1489, T1485, T1562.001)\n\n- **Phishing to RAT/ransomware (Marriott typo, Amnesia RAT)**\n  - Log sources: Secure email gateway, EDR, DNS/proxy, browser logs.  \n  - Focus: messages from newly registered domains, brand impersonation, Office/PDF → script → network connection patterns, new persistence (scheduled tasks, Run keys) after email clicks. (T1566, T1204, T1059, T1053)\n\n- **SSO / IdP abuse (ShinyHunters claims)**\n  - Log sources: IdP/SSO (Okta, Entra ID, etc.), VPN, CASB, SaaS audit logs.  \n  - Focus: new devices/locations for admins, failed MFA followed by success from different IP, new OAuth apps with high‑priv scopes, unusual token refresh patterns. (T1078, T1550, T1552.004)\n\n- **Cloud/SaaS account data theft**\n  - Log sources: Cloud storage (O365, GDrive, S3), data egress gateways, proxy.  \n  - Focus: large downloads/exports from SSO-integrated apps, API key usage from atypical regions, mass access to user lists or profile endpoints. (T1530, T1114, T1567)\n\n- **Potential AI‑assisted / LLM prompt‑injection abuse (early trend)**\n  - Log sources: Any internal AI/LLM gateway logs, application logs using LLMs for automation.  \n  - Focus: prompts containing “ignore previous instructions,” attempts to exfiltrate secrets/configs, abnormal automated actions triggered by LLM outputs. (T1608.005, T1552)\n\n- **General credential & encryption‑related issues (BitLocker / key handling)**\n  - Log sources: Endpoint management, BitLocker logs, HSM/KMS logs.  \n  - Focus: unexpected BitLocker recovery key access, mass key export attempts, changes in disk encryption policies on user or server segments. (T1555, T1552)\n\n---\n\n### Medium-term follow-ups\n\n- **Urgent VMware vCenter patch & hardening program**\n  - Coordinate with infrastructure to patch all affected vCenter builds; restrict management planes to admin networks/VPN; review role assignments and service accounts.\n\n- **BYOVD exposure assessment**\n  - Work with endpoint/security engineering to ensure kernel driver blocklists are enabled (e.g., Windows vulnerable driver blocklist, EDR custom blocklists); remove obsolete drivers commonly abused in BYOVD attacks.\n\n- **SSO / IdP security posture review**\n  - Enforce strong MFA for all accounts (especially admin); implement risk‑based and device‑based conditional access; review and prune unused SAML/OIDC integrations and OAuth consents.\n\n- **Phishing resilience & user awareness**\n  - Refresh training with focus on brand impersonation and subtle domain typos; reinforce “report suspicious email” process; run a targeted simulation around lookalike domains.\n\n- **Backup, recovery, and wiper‑resilience testing**\n  - Validate RPO/RTO and test restores for key systems; ensure at least one offline/immutable backup tier; document and rehearse response runbooks for destructive attacks and mass encryption.\n\n- **Third‑party and data‑sharing risk (in light of breaches & key disclosure news)**\n  - Review critical vendors handling customer data and encryption keys; confirm contractual and technical controls on key management, law enforcement access, and incident notification SLAs.",
  "highlights": [
    {
      "title": "Hackers Use ‘rn’ Typo Trick to Impersonate Marriott in New Phishing Attack",
      "link": "https://cybersecuritynews.com/rn-typo-phishing-attack/",
      "source": "Cyber Security News (FeedBurner)",
      "published": "2026-01-25T04:06:40+00:00",
      "smart_groups": [
        "Curated",
        "Exploit / PoC",
        "Phishing / Social Engineering"
      ],
      "curated": true
    },
    {
      "title": "GitHub - dereeqw/BlackBerryC2: Encrypted command‑and‑control (C2) research framework for cybersecurity education, red team labs, and secure client‑server communication experiments.",
      "link": "https://www.reddit.com/r/cybersecurity/comments/1qm5xib/github_dereeqwblackberryc2_encrypted/",
      "source": "/r/cybersecurity",
      "published": "2026-01-25T01:56:33+00:00",
      "smart_groups": [
        "Curated",
        "Vulnerabilities / CVEs"
      ],
      "curated": true
    },
    {
      "title": "BREAKMEIFYOUCAN! - Exploiting Keyspace Reduction and Relay Attacks in 3DES and AES-protected NFC Technologies",
      "link": "https://www.reddit.com/r/netsec/comments/1qm52ob/breakmeifyoucan_exploiting_keyspace_reduction_and/",
      "source": "/r/netsec",
      "published": "2026-01-25T01:18:14+00:00",
      "smart_groups": [
        "Curated",
        "Exploit / PoC"
      ],
      "curated": true
    },
    {
      "title": "Security Affairs newsletter Round 560 by Pierluigi Paganini – INTERNATIONAL EDITION",
      "link": "https://securityaffairs.com/187290/breaking-news/security-affairs-newsletter-round-560-by-pierluigi-paganini-international-edition.html",
      "source": "Security Affairs",
      "published": "2026-01-25T00:29:16+00:00",
      "smart_groups": [
        "Curated",
        "Ransomware"
      ],
      "curated": true
    },
    {
      "title": "ICE Asks Companies About ‘Ad Tech and Big Data’ Tools It Could Use in Investigations",
      "link": "https://www.wired.com/story/ice-asks-companies-about-ad-tech-and-big-data-tools/",
      "source": "Security News: Cybersecurity, Hacks, Privacy, National Security",
      "published": "2026-01-24T22:14:57+00:00",
      "smart_groups": [
        "Curated",
        "Vulnerabilities / CVEs"
      ],
      "curated": true
    },
    {
      "title": "Swipe, Plug-in, Pwned: Researchers Find New Ways to Hack Vehicles",
      "link": "https://www.reddit.com/r/cybersecurity/comments/1qlwsen/swipe_plugin_pwned_researchers_find_new_ways_to/",
      "source": "/r/cybersecurity",
      "published": "2026-01-24T19:45:50+00:00",
      "smart_groups": [
        "Curated",
        "Vulnerabilities / CVEs",
        "Exploit / PoC"
      ],
      "curated": true
    },
    {
      "title": "CISA confirms active exploitation of four enterprise software bugs",
      "link": "https://www.reddit.com/r/cybersecurity/comments/1qlwoqv/cisa_confirms_active_exploitation_of_four/",
      "source": "/r/cybersecurity",
      "published": "2026-01-24T19:42:03+00:00",
      "smart_groups": [
        "Curated",
        "Exploit / PoC"
      ],
      "curated": true
    },
    {
      "title": "Osiris ransomware emerges, leveraging BYOVD technique to kill security tools",
      "link": "https://securityaffairs.com/187279/security/osiris-ransomware-emerges-leveraging-byovd-technique-to-kill-security-tools.html",
      "source": "Security Affairs",
      "published": "2026-01-24T18:17:13+00:00",
      "smart_groups": [
        "Curated",
        "Ransomware"
      ],
      "curated": true
    },
    {
      "title": "RustyPP: A C++20 library and Clang tool to enforce Rust-like safety and mutability.",
      "link": "https://www.reddit.com/r/cybersecurity/comments/1qlucz4/rustypp_a_c20_library_and_clang_tool_to_enforce/",
      "source": "/r/cybersecurity",
      "published": "2026-01-24T18:17:06+00:00",
      "smart_groups": [
        "Curated",
        "Vulnerabilities / CVEs",
        "Windows / Microsoft",
        "Linux / Unix"
      ],
      "curated": true
    },
    {
      "title": "After 4 years managing an NDR product, I've realized something is fundamentally broken",
      "link": "https://www.reddit.com/r/cybersecurity/comments/1qlsisf/after_4_years_managing_an_ndr_product_ive/",
      "source": "/r/cybersecurity",
      "published": "2026-01-24T17:09:26+00:00",
      "smart_groups": [
        "Curated",
        "Vulnerabilities / CVEs"
      ],
      "curated": true
    }
  ]
}